# scripts/systemd/console-sync.service
[Unit]
Description=yQuant Stock Master Data Sync (Console)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/srv/yquant/console
ExecStart=/usr/bin/dotnet /srv/yquant/console/yQuant.App.Console.dll catalog auto
TimeoutStartSec=1800
