# scripts/systemd/console-sync.timer
[Unit]
Description=yQuant Stock Master Data Sync Timer (Auto Check every 5min)
Requires=console-sync.service

[Timer]
# 5분마다 실행하여 각 국가별 개장 10분 전 타이밍 체크
OnCalendar=*:0/5
# 시스템 부팅 후 10분 뒤에도 실행 (놓친 실행 보완)
OnBootSec=10min
# 타이머가 놓친 실행을 즉시 수행
Persistent=true

[Install]
WantedBy=timers.target
